# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/bin/bash -c "mongod"

[program:tomcat]
command=/bin/bash -c "catalina.sh run"

